# btrfs-keyring.service
[Service]
ExecStart=/usr/bin/gnome-keyring-daemon -c ssh,secrets -f --unlock
